# assets/bash/profile.d/40-man.sh
export MANPATH="$MANPATH:$HOME/.nidus_site_config/nidus/assets/man"
