# assets/bash/profile.d/45-exec.sh
export PATH=$HOME/.nidus_site_config/exec:$PATH
export PATH=$HOME/.nidus_site_config/nidus/assets/exec:$PATH
