# assets/bash/profile.d/45-home-exec.sh
export PATH=$HOME/.nidus_site_config/nidus/assets/exec:$PATH
