# assets/bash/profile.d/90-nidus-exec.sh
export PATH=$HOME/.nidus_site_config/exec:$PATH
