# assets/bash/profile.d/99-banner.sh
if [ "$__NIDUS_BANNER_ENABLED" = yes ]; then
    cat $HOME/.nidus_site_config/nidus/assets/bash/banner
fi
